# systemd/dns-dashboard.service
[Unit]
Description=DNS Analytics Dashboard
After=network.target clickhouse-server.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/dns-dashboard
ExecStart=/opt/dns-dashboard/dns-dashboard
Restart=always
RestartSec=5s
Environment="CLICKHOUSE_DSN=tcp://127.0.0.1:9000?database=dns"
Environment="LISTEN_ADDR=:8080"
Environment="DASHBOARD_USER=admin"
Environment="DASHBOARD_PASS=change_me"
# Ensure simple file descriptor limits are high enough
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
